# deploy/linux/chat-app.service
[Unit]
Description=Realtime Chat App
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/realtime-chat-app
EnvironmentFile=/opt/realtime-chat-app/backend/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
